# Documentation/admin-guide/deepsleep.service
# Example systemd unit for applying deepsleep settings at boot
# Install: copy to /etc/systemd/system/deepsleep.service
[Unit]
Description=Apply deep-sleep friendly runtime power settings
After=multi-user.target
ConditionPathExists=/sys

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/deepsleep-apply.sh enable
ExecStop=/usr/local/sbin/deepsleep-apply.sh disable
TimeoutSec=30

[Install]
WantedBy=multi-user.target
