# Documentation/admin-guide/gaming-mode.service
# Example systemd unit for applying gaming-mode at boot
# Install: copy to /etc/systemd/system/gaming-mode.service

[Unit]
Description=Apply Gaming Mode sysctl settings
After=network.target
ConditionPathExists=/proc/sys/gaming/gaming_mode

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/gaming-mode-apply.sh enable
ExecStop=/usr/local/sbin/gaming-mode-apply.sh disable
TimeoutSec=30

[Install]
WantedBy=multi-user.target
